# layout/DEBIAN/postrm
#!/bin/sh

launchctl unload /Library/LaunchDaemons/com.apple.splashboardd.plist
launchctl load /Library/LaunchDaemons/com.apple.splashboardd.plist
